# deploy/bonfire-rewe-scrape.timer
[Unit]
Description=Bonfire - poll GMX for new REWE eBons daily

[Timer]
# Daily, and shortly after boot. Catches up if the Pi was off (Persistent).
# Groceries arrive once or twice a week; the dashboard's "Fetch mails"
# button covers on-demand fetches in between.
OnBootSec=3min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
